# src/ctp_server.service
[Unit]
Description=Start CTP server.

[Service]
Type=simple
ExecStart=/bin/bash /home/ctp/ctp_server.sh

[Install]
WantedBy=multi-user.target